# dotfiles/systemd/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin renan --noclear %I $TERM
